# src/config/celery.conf
[program:celer]
command=/home/arslan_ushefa/venv/bin/celery worker -A culinary_project --loglevel=INFO
directory=/home/arslan_ushefa/culinary_project
user=arslan_ushefa
autostart=true
autorestart=true
stderr_logfile=/home/arslan_ushefa/culinary_project/logs/celeryd.log
redirect_stderr=true
startsecs=10

#; Need to wait for currently executing tasks to finish at shutdown.
#; Increase this if you have very long running tasks.
stopwaitsecs = 600

#; When resorting to send SIGKILL to the program to terminate it
#; send SIGKILL to its whole process group instead,
#taking care of its children as well.
killasgroup=true

#; if rabbitmq is supervised, set its priority higher
#; so it starts first
priority=998